# decent-hermes/hermes-patch/com.decenthermes.patch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.decenthermes.patch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/claw/.hermes/decent-hermes-patch/apply.sh</string>
    </array>

    <!-- Fire when hermes-agent's main branch ref changes (i.e. git pull landed) -->
    <key>WatchPaths</key>
    <array>
        <string>/Users/claw/.hermes/hermes-agent/.git/refs/heads/main</string>
    </array>

    <!-- Also run once at login, in case an update happened while machine was off -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/claw/.hermes/logs/decentchat-patch.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/claw/.hermes/logs/decentchat-patch.log</string>

    <!-- Throttle: don't re-run more than once per 30s -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
