# decent-hermes/hermes-patch/com.decenthermes.patch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  decent-hermes patch watcher — re-applies the DecentChat adapter patch
  to hermes-agent after any git update.

  This file is a TEMPLATE. The postinstall.cjs script substitutes
  __HOME__ → $HOME when writing the resolved copy to ~/Library/LaunchAgents/.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.decenthermes.patch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/.hermes/decent-hermes-patch/apply.sh</string>
    </array>

    <!-- Fire when hermes-agent's main branch ref changes (git pull landed new commits) -->
    <key>WatchPaths</key>
    <array>
        <string>__HOME__/.hermes/hermes-agent/.git/refs/heads/main</string>
    </array>

    <!-- Also run once at login, in case an update happened while machine was off -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.hermes/logs/decentchat-patch.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.hermes/logs/decentchat-patch.log</string>

    <!-- Throttle: don't re-run more than once per 30s -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
